# agentplatform/deploy/agentplatform-api.service
[Unit]
Description=PitronAgent FastAPI
After=network.target

[Service]
Type=exec
User=www-data
WorkingDirectory=/var/www/PitronAI_AgenticAI/agentplatform
EnvironmentFile=/var/www/PitronAI_AgenticAI/agentplatform/.env
ExecStart=/var/www/PitronAI_AgenticAI/agentplatform/venv/bin/uvicorn api.main:app \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --log-level info
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agentplatform-api

[Install]
WantedBy=multi-user.target
